# Mercury/X18G/squashfs-root/usr/lib/opkg/info/dms.list
/conf/model.xml
/conf/template/upnp.json
/bin/dms
/conf/template/cloud_config.json
/bin/dmclient
/conf/template/wds.json
/etc/init.d/dms
/etc/init.d/wlan
/conf/template/dhcpd.json
/conf/template/wan_port_detect.json
/etc/mcbDesc.conf
/etc/plugin_config.json
/conf/template/plugin_config.json
/conf/template/uhttpd.json
/conf/template/wireless.json
/conf/template/cloud_status.json
/conf/template/network.json
/conf/template/firewall.json
/conf/cfgSyncBlackList.xml
/conf/template/protocol.json
/conf/template/data_template_test.json
/conf/template/system.json
/bin/smp.sh
