# Xiaomi/AX5400/squashfs-root/usr/lib/opkg/info/dualwan.list
/etc/init.d/dualwan
/usr/sbin/dualwan.sh
/etc/config/dualwan
/usr/sbin/dualwan.lua
/usr/sbin/dualwan.util.sh
