# Xiaomi/AX5400/squashfs-root/usr/lib/opkg/info/kmod-shortcut-fe.list
/lib/modules/4.4.60/shortcut-fe-ipv6.ko
/usr/bin/sfe_dump
/lib/modules/4.4.60/shortcut-fe.ko
/etc/init.d/shortcut-fe
