# Xiaomi/AX5400/squashfs-root/usr/lib/opkg/info/miled.list
/etc/config/miled
/lib/config_post_ota/miled_cfg_update.sh
/usr/sbin/miled
/etc/init.d/miled
