# Xiaomi/AX5400/squashfs-root/usr/lib/opkg/info/miwifi-dualwifi.list
/lib/config_post_ota/dualwifi_cfg_update.sh
/etc/init.d/dualwifi
/usr/sbin/dualwifi.sh
/etc/config/dualwifi
