# Xiaomi/AX5400/squashfs-root/usr/lib/opkg/info/statpoints_c.list
/usr/sbin/statpoints
/usr/sbin/mtd_crash_log
/usr/sbin/ecos_pair_verify
/etc/init.d/statpoints
