# Xiaomi/AX5400/squashfs-root/usr/lib/opkg/info/xqled.list
/usr/sbin/xqled_spi
/usr/sbin/xqled
/etc/init.d/xqled
/lib/xqled/leds_spi.sh
