# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/client_monitor.list
/usr/sbin/client_monitor
/etc/init.d/client_monitor
