# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/cs-base-files.list
/etc/config/meshinfo
/etc/config/csfw
/etc/init.d/cs_boot
/etc/hotplug.d/switch/001-switch.sh
/etc/config/wireless
/etc/uci-defaults/13_cs-default
/etc/config/network
