# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/ddns.list
/etc/init.d/noip2
/etc/init.d/yddns
/etc/init.d/orayddns
/usr/bin/orayddns
/etc/config/ddns
/usr/bin/noip2
/usr/bin/yddns
