# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/delay_effect.list
/etc/init.d/apply
/usr/bin/apply
/usr/bin/lktos_reload
