# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/detect_internet.list
/usr/bin/detect_internet
/etc/init.d/detect_internet
