# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/luci-base.list
/usr/lib/lua/luci/view/view.htm
/usr/lib/lua/luci/store.lua
/etc/config/luci
/usr/lib/lua/luci/template.lua
/etc/init.d/ucitrack
/usr/lib/lua/luci/view/header.htm
/usr/lib/lua/luci/view/empty_node_placeholder.htm
/sbin/luci-reload
/usr/lib/lua/luci/view/footer.htm
/usr/lib/lua/luci/controller/admin/uci.lua
/usr/lib/lua/luci/view/error404.htm
/usr/lib/lua/luci/sys/zoneinfo.lua
/usr/lib/lua/luci/config.lua
/usr/lib/lua/luci/sys.lua
/usr/lib/lua/luci/view/csrftoken.htm
/etc/config/ucitrack
/usr/lib/lua/luci/sys/zoneinfo/tzoffset.lua
/usr/lib/lua/luci/sys/zoneinfo/tzdata.lua
/usr/share/acl.d/luci-base.json
/usr/lib/lua/luci/controller/admin/index.lua
/lib/upgrade/keep.d/luci-base
/usr/lib/lua/luci/ccache.lua
/usr/lib/lua/luci/xml.lua
/usr/lib/lua/luci/model/uci.lua
/usr/lib/lua/luci/cacheloader.lua
/www/cgi-bin/luci
/usr/share/luci/menu.d/luci-base.json
/etc/luci-uploads/.placeholder
/usr/lib/lua/luci/i18n.lua
/usr/lib/lua/luci/version.lua
/usr/lib/lua/luci/view/sysauth.htm
/usr/lib/lua/luci/dispatcher.lua
/usr/libexec/rpcd/luci
/usr/share/rpcd/acl.d/luci-base.json
/usr/lib/lua/luci/sgi/uhttpd.lua
/usr/lib/lua/luci/template/parser.so
/usr/lib/lua/luci/view/error500.htm
/www/index.html
/usr/lib/lua/luci/view/indexer.htm
/usr/lib/lua/luci/sgi/cgi.lua
