# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/netcwmp.list
/usr/sbin/config-cwmp-dhcpopt.sh
/etc/init.d/cwmp
/usr/sbin/cwmpd
