# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/omcproxy.list
/etc/config/omcproxy
/usr/sbin/omcproxy
/etc/init.d/omcproxy
