# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/shttpd.list
/etc/ssl/certs/server.pem
/usr/sbin/shttpd
/etc/init.d/shttpd
