# TOTOLINK/X6000R/squashfs-root/usr/lib/opkg/info/statistics.list
/usr/bin/cs_statistics
/etc/init.d/cs_statistics
