# Mercury/A15/v1.0.6/squashfs-root/etc/hotplug.d/usb/10-usb
#!/bin/sh

# Copyright (C) 2009 OpenWrt.org

PATH=$PATH:/tmp/usr/bin

case "$ACTION" in
	add)
		gpio_ctrl --usb_insert
		# update LEDs
		;;
	remove)
		# update LEDs
		;;
esac	
		
