# Mercury/X30G/squashfs-root/usr/lib/opkg/info/kmod-qca-wifi-flash_16mb-profile.list
/lib/wifi/qcacommands_ol_radio.xml
/ini/internal/QCA5018_i.ini
/ini/QCA5018.ini
/usr/sbin/thermaltool
/lib/firmware/IPQ5018/firmware_rdp_feature.ini
/ini/QCN6122.ini
/usr/lib/libwlanfw_wmi.so
/lib/firmware/qcn6122/m3_fw.b00
/lib/firmware/qcn6122/m3_fw.b01
/etc/init.d/wireless_start
/lib/firmware/qcn6122/m3_fw.b02
/lib/modules/4.4.60/ath_pktlog.ko
/etc/config/qcacfg80211
/usr/sbin/apstats
/lib/firmware/IPQ5018/m3_fw.b00
/lib/firmware/IPQ5018/m3_fw.b01
/lib/firmware/IPQ5018/m3_fw.b02
/lib/modules/4.4.60/qdf.ko
/lib/firmware/qcn6122/bdwlan.b50
/lib/wifi/qca_ftm_init.sh
/usr/sbin/spectraltool
/lib/wifi/tools_config
/lib/modules/4.4.60/qca_ol.ko
/ini/internal/global_i.ini
/lib/firmware/IPQ5018/q6_fw.mdt
/usr/sbin/cfg80211tool
/lib/firmware/IPQ5018/m3_fw.flist
/usr/lib/libqca_tools.so
/usr/sbin/radartool
/lib/modules/4.4.60/wifi_3_0.ko
/lib/wifi/qcawifi_countrycode.txt
/etc/init.d/boot-ftm
/usr/sbin/wifistats
/lib/firmware/IPQ5018/q6_fw.flist
/usr/sbin/wifitool
/lib/firmware/IPQ5018/bdwlan.b10
/lib/firmware/qcn6122/qdss_trace_config.bin
/lib/firmware/qcn6122/m3_fw.mdt
/lib/modules/4.4.60/asf.ko
/lib/firmware/qcn6122/m3_fw.flist
/lib/wifi/qcawificfg80211.sh
/lib/firmware/IPQ5018/m3_fw.mdt
/usr/lib/libwlanfw_htt.so
/lib/wifi/qca-wifi-modules
/etc/init.d/diag_socket_app
/ini/internal/QCN6122_i.ini
/usr/sbin/wlanconfig
/lib/wifi/qcacommands_vap.xml
/lib/firmware/IPQ5018/qdss_trace_config.bin
/lib/firmware/IPQ5018/fw_version.txt
/ini/global.ini
/lib/firmware/IPQ5018/firmware_rdp_feature_512P.ini
/ini/QCA5018_TP.ini
/ini/QCN6122_TP.ini
/lib/modules/4.4.60/umac.ko
/lib/modules/4.4.60/qca_spectral.ko
/lib/firmware/IPQ5018/q6_fw.b00
/lib/firmware/IPQ5018/q6_fw.b01
/lib/firmware/IPQ5018/q6_fw.b02
/lib/firmware/IPQ5018/q6_fw.b03
/lib/firmware/IPQ5018/q6_fw.b04
/lib/firmware/IPQ5018/q6_fw.b05
/lib/firmware/IPQ5018/q6_fw.b07
/lib/firmware/IPQ5018/q6_fw.b08
/lib/firmware/IPQ5018/q6_fw.b09
/usr/sbin/pktlogconf
/lib/firmware/IPQ5018/q6_fw.b10
/lib/firmware/IPQ5018/q6_fw.b11
/lib/firmware/IPQ5018/q6_fw.b13
/lib/firmware/IPQ5018/q6_fw.b14
/lib/firmware/IPQ5018/q6_fw.b15
/lib/firmware/IPQ5018/q6_fw.b16
/lib/firmware/IPQ5018/q6_fw.b17
/lib/firmware/IPQ5018/q6_fw.b18
/lib/modules/4.4.60/mem_manager.ko
/lib/firmware/IPQ5018/q6_fw.b20
/lib/firmware/IPQ5018/q6_fw.b21
/lib/firmware/IPQ5018/q6_fw.b22
/lib/firmware/IPQ5018/q6_fw.b23
/lib/wifi/qcawifi.sh
/lib/firmware/IPQ5018/q6_fw.b25
/lib/firmware/IPQ5018/q6_fw.b26
