# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/baidupan.list
/usr/sbin/baidupan_uci_lock.sh
/usr/sbin/baidupan.lua
/usr/sbin/kill_baidupan.sh
/etc/init.d/baidupan
/etc/config/baidupan
/usr/sbin/baidupan
