# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/datacenter.list
/usr/share/datacenter/appdata/notuninstallable/app_package/2882303761517281003
/usr/share/datacenter/pw
/opt/plugin_sandbox/dropbear.init
/opt/plugin_sandbox/passwd
/usr/share/datacenter/cp_preinstall_plugins_impl.sh
/usr/sbin/pluginControllor
/usr/sbin/pushClient
/etc/datacenterconfig/localization/en.txt
/usr/sbin/plugincenter
/opt/plugin_sandbox/dropbear
/usr/share/datacenter/sdk_level
/opt/plugin_sandbox/create_sandbox.sh
/usr/lib/libxmrouter.so
/etc/datacenterconfig/localization/hk.txt
/usr/sbin/plugin_stop_impl_standalone_hd.sh
/usr/sbin/indexservice
/usr/sbin/chroot_umout.sh
/usr/sbin/installplugin
/etc/init.d/indexservice.init
/opt/plugin_sandbox/start_dropbear.sh
/usr/sbin/mount_download_folder.sh
/usr/bin/chroot
/usr/share/datacenter/appdata/notuninstallable/app_infos/2882303761517280984.manifest
/opt/plugin_sandbox/prepare_sandbox.sh
/etc/init.d/datacenter
/usr/sbin/plugin_start_impl_standalone_hd.sh
/usr/sbin/prepare_mount_dev.sh
/etc/init.d/plugincenter
/usr/sbin/prepare_chroot.sh
/opt/plugin_sandbox/shadow
/usr/lib/libxmrouter.so.1.0.0
/etc/clear.sh
/usr/sbin/datacenterClient
/etc/datacenterconfig/localization/kr.txt
/etc/init.d/plugin_start_script.sh
/usr/sbin/shareUpdate
/etc/datacenterconfig/localization/tw.txt
/opt/plugin_sandbox/sandbox_umount.sh
/usr/sbin/copy_plugin_chroot_file
/usr/share/datacenter/文件共享说明.jpg
/usr/sbin/mount_userdata.sh
/usr/sbin/prepare_supervisord.sh
/usr/share/datacenter/appdata/notuninstallable/app_package/2882303761517280984
/opt/plugin_sandbox/group
/usr/sbin/prepare_mount_cache_dir.sh
/opt/plugin_sandbox/stop_dropbear.sh
/usr/lib/libxmrouter.so.1
/etc/init.d/cp_preinstall_plugins.sh
/etc/datacenterconfig/localization/cn.txt
/usr/sbin/datacenter
/opt/plugin_sandbox/sandbox_white_list.conf
/usr/lib/libroutermain.a
/etc/datacenterconfig/localization/country.txt
/usr/share/datacenter/init.template
/usr/sbin/prepare_mount_proc.sh
/usr/sbin/chroot_run_script.sh
