# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/docker.list
/usr/sbin/mi_docker_disk_manage.sh
/etc/config/mi_docker
/etc/init.d/mi_docker
