# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/dualwan.list
/etc/init.d/dualwan
/usr/sbin/dualwan.firewall
/etc/config/dualwan
/usr/sbin/dualwan.lua
