# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/filetunnel.list
/etc/init.d/stunserver
/opt/filetunnel/stunserver
/etc/init.d/filetunnel
/opt/filetunnel/test_file_1
/opt/filetunnel/tunnelserver
/opt/filetunnel/tunneltool
/opt/filetunnel/test_file
