# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/flashtest.list
/etc/init.d/flashtest
/usr/sbin/flashtest
