# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/kmod-qca-nss-drv-igs.list
/lib/modules/4.4.60/act_nssmirred.ko
/etc/init.d/qca-nss-mirred
