# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/kmod-qca-nss-drv-ovpn-mgr.list
/etc/init.d/qca-nss-ovpn
/lib/modules/4.4.60/qca-nss-ovpn-mgr.ko
