# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/kmod-qca-nss-drv-ipsecmgr-klips.list
/etc/init.d/qca-nss-ipsec
/lib/modules/4.4.60/qca-nss-ipsec-klips.ko
