# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/makeswap.list
/etc/config/makeswap
/etc/init.d/makeswap
