# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/miled.list
/etc/config/miled
/usr/sbin/miled
/etc/init.d/miled
