# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/qca-hostap-macsec.list
/etc/init.d/qca-hostapd
/lib/wifi/hostapd.sh
/usr/sbin/hostapd-macsec
