# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/qca-wpa-supplicant-macsec.list
/lib/wifi/wpa_supplicant.sh
/usr/sbin/wpa_supplicant-macsec
/etc/init.d/macsec
/etc/config/macsec
/etc/init.d/qca-wpa-supplicant
