# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/samba4-server.list
/usr/bin/smbcontrol
/etc/config/samba
/usr/bin/smbstatus
/usr/bin/smbtree
/usr/sbin/nmbd
/usr/bin/smbcacls
/usr/bin/testparm
/etc/samba/smb.conf.template
/etc/samba/valid.dat
/etc/init.d/samba
/etc/samba/secrets.tdb
/usr/bin/smbcquotas
/etc/samba/upcase.dat
/usr/sbin/smbpasswd
/usr/sbin/smbd
/usr/sbin/samba_multicall
/usr/bin/nmblookup
/etc/samba/schannel_store.tdb
/etc/openldap/schema/samba.schema
/etc/samba/lowcase.dat
