# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/sysapihttpd.list
/etc/sysapihttpd/htdocs/miwifi.html
/etc/mihttpd/win-utf
/usr/sbin/sysapihttpd
/etc/sysapihttpd/nginx-ssl.key
/etc/mihttpd/koi-utf
/usr/sbin/mihttpd
/etc/mihttpd/mihttpd.conf
/lib/upgrade/keep.d/sysapihttpd
/etc/sysapihttpd/miwifi-webinitrd-https.conf
/etc/sysapihttpd/mime.types
/etc/mihttpd/htdocs/index.html
/etc/mihttpd/mime.types
/etc/sysapihttpd/cert.crt
/etc/mihttpd/htdocs/miwifi.html
/etc/sysapihttpd/koi-win
/etc/init.d/sysapihttpd
/etc/sysapihttpd/cert.key
/etc/sysapihttpd/win-utf
/etc/sysapihttpd/koi-utf
/etc/mihttpd/fastcgi_params
/etc/sysapihttpd/init
/etc/sysapihttpd/htdocs/index.html
/etc/sysapihttpd/stat.conf
/etc/sysapihttpd/fastcgi-proxy-tcp.conf
/etc/sysapihttpd/htdocs/wifishare.html
/etc/mihttpd/miwifi-webinitrd-mihttpd.conf
/etc/mihttpd/koi-win
/etc/sysapihttpd/fastcgi_params
/etc/sysapihttpd/sysapihttpd.conf
/etc/mihttpd/fastcgi-proxy-tcp-mihttpd.conf
/etc/sysapihttpd/miwifi-webinitrd.conf
/etc/sysapihttpd/nginx-ssl.crt
/etc/init.d/mihttpd
