# Xiaomi/AX9000/squashfs-root/usr/lib64/opkg/info/xq_info_sync_mqtt.list
/sbin/whc_to_re_common_api.sh
/sbin/collect_topo_info.sh
/usr/sbin/xq_info_sync_mqtt
/sbin/get_samba_status.sh
/sbin/samba_status_change.sh
/sbin/xqwhc_backhaul_maclist_reload.sh
/sbin/miled_ctrl.sh
/sbin/bind_client.sh
/sbin/xqwhc_channel.sh
/sbin/cap_whc_common_push.sh
/etc/config/bind
/sbin/xqwhc_gw_update.sh
/sbin/cap_push_backhaul_whitelist.sh
/sbin/re_do_action_common_api.sh
/etc/init.d/xq_info_sync_mqtt
/sbin/xqwhc_maclist_reload.sh
/sbin/sync_tz.lua
/sbin/xqwhc_info_sync.sh
/sbin/re_restore.sh
