# LB-LINK/X26/X26(AC12)_V1.0.1_TFTP/squashfs-root/etc/crontabs/root
*/5 * * * * /bin/set_timezone.sh 
*/5 * * * * /bin/check_progress_status.sh
*/30 * * * * /bin/check_storage.sh
