# Mercury/A15/v1.0.6/squashfs-root/usr/lib/opkg/info/rtl_hostapd.list
/etc/init.d/wlan_hostapd
/conf/hostapd_5g.conf
/usr/sbin/hostapd
/usr/sbin/rtl_ftm_init.sh
/usr/sbin/hostapd_cli
/conf/hostapd_2g.conf
