# Mercury/A15/v1.0.6/squashfs-root/usr/lib/opkg/info/rtl_wpa_supplicant.list
/conf/wpa_supplicant.conf
/etc/init.d/wpas_start
/usr/sbin/wpa_supplicant
/usr/sbin/wpa_cli
