# LB-LINK/X26/X26(AC12)_V1.0.1_TFTP/squashfs-root/usr/lib/opkg/info/base-files.list
/usr/sbin/mac
/etc/openwrt_version
/bin/check_progress_status.sh
/etc/banner.failsafe
/etc/resolv.conf
/etc/mac_alias_list
/etc/init.d/done
/bin/config_generate
/var
/sbin/wps_led.sh
/etc/opkg/keys/de98a2dd1d0f8a07
/bin/wifi_wps.sh
/etc/init.d/system
/etc/preinit
/etc/shells
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/device_info
/sbin/get_7628_lan_wan_mac.sh
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/bin/restart_firewall.sh
/lib/preinit/10_indicate_preinit
/bin/enable_ipv6.sh
/bin/restart_easycwmpd.sh
/etc/blacklist
/lib/functions/uci-defaults-new.sh
/etc/config/system
/bin/check_bridge.sh
/sbin/green_led_blink.sh
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/etc/terminal_lock_file
/etc/opkg/keys/53bad1233d4c98c5
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/wifi_schedule_list
/etc/banner
/bin/restart_dnsmasq.sh
/sbin/platform_pre_upgrade
/etc/shadow
/etc/qos_limit_list
/etc/rc.common
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/10_indicate_failsafe
/init
/sbin/reset_lan_port.sh
/etc/udhcpd.config
/lib/preinit/80_mount_root
/sbin/backup_lan.sh
/etc/rc.local
/lib/functions/uci-defaults.sh
/sbin/red_led_blink.sh
/etc/init.d/led
/etc/config/internet
/sbin/config_nas.sh
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/network_led.sh
/sbin/hotplug-call
/etc/crontabs/root
/bin/restart_network.sh
/etc/rc.button/failsafe
/etc/parental_list
/etc/services
/lib/functions/system.sh
/bin/restart_ddns.sh
/bin/check_storage.sh
/bin/board_detect
/lib/upgrade/common.sh
/etc/rc.button/power
/etc/sysupgrade.conf
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/periodic_restart_list
/etc/rc.button/reset
/sbin/wifi_connect.sh
/lib/preinit/30_failsafe_wait
/bin/set_timezone.sh
/etc/init.d/telnet
/etc/timezone_file
/etc/opkg/keys/758078c5a78d3618
/sbin/firstboot
/sbin/wifi_setup.sh
/bin/login.sh
/sbin/reset_apcli_status.sh
/etc/hotplug.d/net/00-sysctl
/sbin/wifi_led_blink_slow.sh
/etc/hosts
/sbin/unset_vlan.sh
/lib/functions.sh
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/sbin/wifi_led_blink_fast.sh
/rom/note
/sbin/qos_rule_edit.sh
/etc/inittab
/sbin/qos.sh
/etc/openwrt_release
/sbin/wifi
/bin/restart_ntpd.sh
/etc/profile
/lib/wifi/mtkwifi.inc
/sbin/mtk_factory_rw.sh
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/TZ
/etc/init.d/umount
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/sbin/set_vlan.sh
/lib/preinit/02_default_set_state
/etc/mtab
/sbin/check_ap.sh
/etc/init.d/boot
/etc/uci-defaults/03_dat2uci
/lib/functions/preinit.sh
/bin/ipcalc.sh
/sbin/wifi.legacy
/sbin/set_wifi_led.sh
/etc/config/iptv
/etc/FW-Version
/etc/port_forward_list
