# LB-LINK/X26/X26(AC12)_V1.0.1_TFTP/squashfs-root/usr/lib/opkg/info/Blink_Cloud.list
/usr/bin/Blink_Cloud
/etc/init.d/Blink_Cloud
