# LB-LINK/X26/X26(AC12)_V1.0.1_TFTP/squashfs-root/usr/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.conf
/etc/dhcp_client_lock_file
/etc/config/dhcp
/etc/hotplug.d/gmac/24-dnsmasq
/etc/hotplug.d/ntp/25-dnsmasqsec
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
