# LB-LINK/X26/X26(AC12)_V1.0.1_TFTP/squashfs-root/usr/lib/opkg/info/sys_daemon.list
/etc/init.d/sys_daemon
/usr/bin/sys_daemon
