# LB-LINK/X26/X26(AC12)_V1.0.1_TFTP/squashfs-root/usr/lib/opkg/info/terminal_list.list
/usr/bin/terminal_list
/etc/init.d/terminal_list
